# duplicate_video_finder/rootfs/etc/services.d/duplicate-video-finder/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Duplicate Video Finder..."

# Change to the app directory
cd /opt/duplicate_video_finder

# Run the Python script
exec python3 /opt/duplicate_video_finder/run.py
