# duplicate-video-finder/rootfs/etc/services.d/duplicate-video-finder/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Duplicate Video Finder..."

# Change to the app directory
cd /app

# Run the Python script
exec python3 /app/run.py
