# duplicate-video-finder/rootfs/etc/services.d/duplicate-video-finder/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Duplicate Video Finder..."

# Change to the app directory
cd /usr/src/app

# Run the application without using exec to prevent PID 1 issues
python3 /usr/src/app/run.py
