# duplicate-video-finder/rootfs/etc/services.d/duplicate-video-finder/run
#!/command/with-contenv bashio

# ==============================================================================
# Start Duplicate Video Finder service
# ==============================================================================
bashio::log.info "Starting Duplicate Video Finder service..."

# Change to the app directory
cd /usr/src/app || bashio::exit.nok "Could not change to app directory"

# Run application (but don't use exec)
python3 /usr/src/app/run.py
