# duplicate-video-finder/rootfs/etc/services.d/duplicate-video-finder/run
#!/bin/bash

# ==============================================================================
# Start Duplicate Video Finder service
# ==============================================================================
echo "Starting Duplicate Video Finder service..."

# Change to the app directory
cd /usr/src/app || exit 1

# Run application (but don't use exec)
python3 run.py
